# init.d/cloudflared
#!/sbin/openrc-run

name=$RC_SVCNAME
description="cloudflared DNS over HTTPS proxy"

depend()
{
	after syslog network-online
}



